# packages/pulp/pulp-maintenance.cron
#!/bin/sh
celery -A pulp.server.async.app call pulp.server.maintenance.monthly.queue_monthly_maintenance
